# configs/crontab
* * * * * sh ~/scripts/git_sync_main.sh
