# configs/crontab
* * * * * sh ~/scripts/git_sync_main.sh
* * * * * sh ~/scripts/git_safe_pull.sh
