# configs/crontab
* * * * * sh ~/scripts/git_sync_branches.sh
* * * * * sh ~/scripts/git_sync_main.sh
