# configs/crontab
*/5 * * * * sh $GIT_PROJECTS_WORKDIR/dotfiles/scripts/git/git_sync_all.sh
*/5 * * * * sh $GIT_PROJECTS_WORKDIR/dotfiles/scripts/git/git_fetch_all.sh
